# systemd.timer
[Unit]
Description=Run /opt/endpoint/backup.sh

[Timer]
OnBootSec=12h
OnUnitActiveSec=1d

[Install]
WantedBy=timers.target
