# systemd.timer
[Unit]
Description=Run /opt/endpoint/backup.sh

[Timer]
OnBootSec=30m
OnUnitActiveSec=1h

[Install]
WantedBy=timers.target
